# jobs-mcp.service
# /etc/systemd/system/jobs-mcp.service
# Run: sudo systemctl enable --now jobs-mcp

[Unit]
Description=Jobs MCP Server (healthcare jobs)
After=network.target mongod.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/home/dev_xpresshealth/jobs-mcp-server
EnvironmentFile=/home/dev_xpresshealth/jobs-mcp-server/.env
ExecStart=/home/dev_xpresshealth/jobs-mcp-server/venv/bin/python jobs_mcp_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
